# mail-mta/sendmail/files/var-spool-mqueue.path


[Path]

DirectoryNotEmpty=/var/spool/mqueue

Unit=sendmail-mta.service

[Install]
WantedBy=default.target

