# net-proxy/squid/files/squid.socket
[Unit]
Description=squid socket 

[Socket]
Transparent=true
ListenStream=0.0.0.0:3128

[Install]
WantedBy=sockets.target
